# software/cansetup.service
[Unit]
Description=CAN Setup

[Service]
ExecStart=/bin/bash -c /home/rover/startupCAN.sh

[Install]
WantedBy=default.target
